# scripts/deploy/apollo_portal.socket
[Unit]
Description=Apollo Portal website socket
[Socket]
ListenStream=/run/gunicorn.sock
[Install]
WantedBy=sockets.target
